# provisioning/ansible/playbooks/app/templates/collector.service.j2
[Service]
ExecStart={{ collector_path }}/bin/start-collector.sh
User={{ ansible_user }}
Group={{ ansible_user }}
