# provisioning/ansible/roles/hadoop/templates/hadoop.service.j2
[Service]
ExecStart={{ hadoop_path }}/sbin/start-all.sh

User={{ hadoop_user }}
Group={{ hadoop_group }}
