# provisioning/ansible/roles/zookeeper/templates/zookeeper.service.j2
[Service]
ExecStart={{ zookeeper_path }}/bin/zkServer.sh start-foreground
ExecStop={{ zookeeper_path }}/bin/zkServer.sh stop
User={{ zookeeper_user }}
Group={{ zookeeper_group }}
