# provisioning/ansible/roles/kafka/templates/kafka.service.j2
[Service]
ExecStart={{ kafka_path }}/bin/kafka-server-start.sh {{ kafka_path }}/config/server.properties
User={{ kafka_user }}
Group={{ kafka_group }}
